# tayga@.service.in
[Unit]
Description=Simple, no-fuss NAT64
After=network.target

[Service]
# To set up an extra tayga service instance, create a new tayga config in
# @sysconfdir@/tayga/instancename.conf.
# Then enable this config as a systemd alias:
#
# systemctl enable tayga@instancename.service

Type=simple
PrivateTmp=true
ExecStart=@sbindir@/tayga -d --config @sysconfdir@/tayga/%i.conf

[Install]
WantedBy=multi-user.target
