# tayga.service.in
[Unit]
Description=Simple, no-fuss NAT64
After=network.target

[Service]
Type=simple
PrivateTmp=true
ExecStart=@sbindir@/tayga -d --config @sysconfdir@/tayga.conf

[Install]
WantedBy=multi-user.target
